# src/etc/systemd/system/demo-deb-package.service
[Unit]
Description=Run package_me.sh

[Service]
ExecStart=/usr/bin/package_me.sh