# src/etc/systemd/system/demo-deb-package.timer
[Unit]
Description=Run package_me.sh every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target